# airanking.service
[Unit]
Description=AIRanking Team Building Game Score System
After=network.target

[Service]
User=www-data
Group=www-data
# 添加补充组 - 确保 www-data 可以访问 jerry 的文件
SupplementaryGroups=jerry
WorkingDirectory=/var/www/airankingx.com
ExecStart=/usr/bin/python3 /var/www/airankingx.com/airankingx.py
Restart=always
RestartSec=10
# 使用独立的日志文件
StandardOutput=journal
StandardError=journal
# 设置环境变量
Environment=PYTHONUNBUFFERED=1
# 确保服务启动时刷新组权限
ExecStartPre=/bin/bash -c 'id www-data | logger -t airanking-service'

[Install]
WantedBy=multi-user.target

